# root/etc/s6-overlay/s6-rc.d/init-mod-jellyfin-introskipper-button/run
#!/usr/bin/with-contenv bash

chown abc /usr/share/jellyfin/web/index.html
